# profile.d/misc.sh
# Default file and dir permission mask
# Subracting from 666 for files and 777 for directories
umask 0002
